# Raspberry-pi/rc.local
#!/bin/sh -e
bash /home/pi/Pi/pi.sh
exit 0